# bot.service
[Unit]
Description=Slurm Cluster Monitor Bot
After=network.target

[Service]
# Replace 'ubuntu' with your actual username on the VPS
User=ubuntu
# Adjust the path to where you cloned the repo
WorkingDirectory=/home/ubuntu/server-notification
# Adjust the path to your python executable (or virtualenv)
ExecStart=/usr/bin/python3 /home/ubuntu/server-notification/monitor.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
